# docker/conf/supervisor/conf.d/php7.4-fpm.conf
[program:php7.4-fpm]
command=/usr/sbin/php-fpm7.4  --nodaemonize --force-stderr --fpm-config /etc/php/7.4/fpm/php-fpm.conf
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
